# roles/smb/templates/smb.automount.j2
[Unit]
Description=nas automount

[Automount]
Where={{ smb_mountpoint }}

[Install]
WantedBy=multi-user.target
