# roles/home/files/systemd-user-units/nas-backup.timer
[Unit]
Description=nas-backup timer

[Timer]
Unit=nas-backup.service
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
