# roles/home/files/systemd-user-units/nas-prune.timer
[Unit]
Description=nas-prune timer

[Timer]
Unit=nas-prune.service
OnCalendar=monthly
Persistent=true

[Install]
WantedBy=timers.target
