# roles/sway/files/systemd-user-units/user-service-template.service
[Unit]
Description=
PartOf=graphical-session.target

[Service]
Type=simple
ExecStart=

[Install]
WantedBy=sway-session.target