# roles/manpack/files/manpack.timer
[Unit]
Description=manpack timer

[Timer]
Unit=manpack.service
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
